# dotfiles/rclone/.config/systemd/user/rclone-sync.timer
[Unit]
Description=Perform an RClone sync periodically
[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target
